# files/containers-root/hd-idle.container
[Unit]
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=hd-idle
Image=ghcr.io/lehmanju/hd-idle:latest
AutoUpdate=registry
PodmanArgs=--privileged
Exec=-i 300 -c ata

[Service]
MaxRuntimeSec=12h
Restart=always

[Install]
WantedBy=multi-user.target default.target