# files/containers/nc_redis.container
[Container]
Pod=nextcloud.pod
ContainerName=nc_redis
AutoUpdate=registry
Image=docker.io/redis:alpine